# contrib/systemd/x0d.socket
[Unit]
Description=x0d socket-based activator

[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target
